# .config/autostart/protonvpn.desktop
[Desktop Entry]
Name=ProtonVPN
Exec=protonvpn-app
Type=Application
