# .config/autostart/telegram.desktop
[Desktop Entry]
Name=Telegram
Exec=telegram-desktop
Type=Application
